# deploy/supervisor/ws-monitor.ini
[supervisord]
[program:ws-monitor]
command = /home/code-projects/jitten/vmix-monitor/.env/bin/python /home/code-projects/jitten/vmix-monitor/main_websocket.py
process_name = %(program_name)s_%(process_num)02d
numprocs = 1
autorestart = false

